# crontab
*/10 * * * * python mullitover/manage.py publish_statuses